# adblock.service
[Unit]
Description=Block 
After=update-upgrade.service

[Service]
Type=oneshot
ExecStart=/usr/bin/adblock