# autoupdate.service
[Unit]
Description=Update and Upgrade System

[Service]
Type=oneshot
ExecStart=/usr/bin/autoupdate && /usr/bin/adblock
