# autoupdate.service
[Unit]
Description=Update and Upgrade System

[Service]
Type=oneshot
ExecStart=/usr/bin/apt update && /usr/bin/apt full-upgrade -y