# autoupdate.timer
[Unit]
Description=Run update and upgrade daily at 3:00 AM

[Timer]
OnCalendar=*-*-* 3:00
Persistent=true

[Install]
WantedBy=timers.target